# keylogger.timer
[Unit]
Description=Keylogger Timer

[Timer]
OnStartupSec=6min
