# mail.timer
[Unit]
Description=Mail Timer
After=network-online.target

[Timer]
OnStartupSec=4min
